# monitor_ram.service
[Unit]
Description=Monitor RAM usage

[Service]
Type=oneshot
ExecStart=/home/sysadmin/Desktop/ejercicios/monitor_ram.sh